# install/projects_overland_server.service
[Unit]
 Description=Incognita Overland Server
 After=multi-user.target

 [Service]
 WorkingDirectory=/home/mnalavadi/incognita
 Type=idle
 ExecStart=/home/mnalavadi/miniconda3/envs/gps/bin/python -m incognita.overland_server
 User=mnalavadi

 [Install]
 WantedBy=multi-user.target
